# scripts/subtitle_pipeline/subtitle-pipeline.timer
[Unit]
Description=Run subtitle pipeline hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=subtitle-pipeline.service

[Install]
WantedBy=timers.target
