# scripts/k3s-state-backup.service
[Unit]
description=backup k3s state

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-k3s-state.sh