# infra/ansible/roles/email-alert/templates/email-alert.service.j2
[Unit]
Description=Send email alert when %i fails

[Service]
Type=oneshot
ExecStart=/usr/local/bin/systemd-email-alert.sh %i